# launchd/com.lx.daily-hot-news.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lx.daily-hot-news</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/lx/Personal/daily-hot-news/.venv/bin/python</string>
        <string>-m</string>
        <string>src.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/lx/Personal/daily-hot-news</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/lx/Personal/daily-hot-news/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/lx/Personal/daily-hot-news/logs/stderr.log</string>
</dict>
</plist>
